# bubbleupnpserver.service
[Unit]
Description=BubbleUPnP Server
Requires=network-online.target
After=network-online.target

[Service]
Type=simple
Restart=on-failure
ExecStart=/usr/bin/bubbleupnpserver

[Install]
WantedBy=multi-user.target
